# packages/vajra-nightly/vajra_nightly-0.1.dev2025060906.tar.gz/vajra_nightly-0.1.dev2025060906/scripts/github_runner_setup/github-runner.service
[Unit]
Description=GitHub Actions Runner GPU Distributor
After=docker.service nvidia-persistenced.service
Requires=docker.service
Documentation=https://github.com/project-vajra/vajra

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/github-runners

# Load the GitHub token as a credential
EnvironmentFile=/opt/github-runners/github-token.env

# Pass the token from the credentials directory to the script
ExecStart=/bin/bash -c "/opt/github-runners/launch_runners.sh /opt/github-runners/runner_spec.json https://github.com/project-vajra/vajra project-vajra ${GITHUB_TOKEN}"
ExecStop=/opt/github-runners/stop_runners.sh

# Restart configuration
Restart=always
RestartSec=30

# Standard output/error logs
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
